# unbound-crypt/services.d/unbound/run
#!/bin/sh

exec /usr/sbin/unbound -d -p -v 
